# gamemaker.service
[Unit]
Description=morganashley	

[Service]
ExecStart=/opt/morganashley/bin/morganashley

[Install]
WantedBy=multi-user.target